# service/sep-server.service
[Unit]
Description=Sep Server
After=network.target

[Service]
ExecStart=/usr/local/bin/sep-server --key {key}
Restart=on-failure
User=nobody

[Install]
WantedBy=multi-user.target
